# failures_server/docker_configs/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/docker/PVmeasure/failures_server/docker_configs/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
